# basrcpd/files/presume-ack.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2020 OpenWrt.org

START=97

start() {
	ip link set can0 down
	ip link set can0 type can presume-ack on
	[ -e /var/run/can2lan.pid ] && \
		ip link set can0 up 
}

stop() {
	ip link set can0 down
	ip link set can0 type can presume-ack off
	[ -e /var/run/can2lan.pid ] && \
		ip link set can0 up 
}

