# can2udp/files/can2lan.init_srse2
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 OpenWrt.org

START=99

start() {
	omega2-ctrl gpiomux set uart1 uart
	slcand -S 500000 /dev/ttyS1 can0
	ifconfig can0 up

	[ -e /www/config/ ] && {
		rm -f /var/run/can2lan.pid
		can2lan -m -c /www/config -b eth0
		# can2lan -m -c /www/config -b apcli0
		# can2lan -m -c /www/config -b br-wlan
	}
}

stop() {
	killall can2lan
}
