# can2udp/files/can2lan.init_slcan_a5v11
#!/bin/sh /etc/rc.common
# Copyright (C) 2017 OpenWrt.org

START=99

start() {
	slcand -S 500000 /dev/ttyUSB0 can0
	# wait a little bit to establish CAN device
	sleep 1
	if [ -e /usr/sbin/ip ]; then
		/usr/sbin/ip link set can0 up type can
	else
		/sbin/ip link set can0 up type can
	fi

	[ -e /www/config/ ] && {
		can2lan -m -c /www/config -b br-lan
	}
}

stop() {
	killall can2lan
}
