# can2udp/files/can2lan.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 OpenWrt.org

START=99

start() {
	modprobe sun4i_can
	ip link set can0 type can bitrate 250000
	ifconfig can0 up
	sleep 10

	[ -e /www/config/ ] && {
		can2lan -m -c /www/config -b br-lan
	}
}

stop() {
	killall can2lan
}
