# can2udp/files/clone-ms2-loco.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2018 OpenWrt.org

START=99

start() {
	# clone every 5 minutes
	# clone-ms2-loco -i 500

	# usining LED on PI14 and push button on PI10	
	clone-ms2-loco -l 270 -p 266
}

stop() {
	killall clone-ms2-loco
}
