# can2udp/files/cs2ping.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2017 OpenWrt.org

START=99

start() {
	# cs2ping -b br-lan
}

stop() {
	killall cs2ping
}
