# can2udp/files/ms2-loco-list.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2020 OpenWrt.org

START=99

start() {
	ms2-loco-list -d
}

stop() {
	killall ms2-loco-list
}
