# can2udp/files/wake-up-links88.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2016 OpenWrt.org

START=99

start() {
	wake-up-links88 -dc "B3=1" -i can0
}

stop() {
	killall wake-up-links88
}
