# ms1relay/files/ms1relay.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2020 OpenWrt.org

START=99

start() {
	ms1relay -d
}

stop() {
	killall ms1relay
}
