# railcontrol/files/railcontrol.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2020 OpenWrt.org

START=99

start() {
	[ -e /opt/railcontrol/ ] && {
		cd /opt/railcontrol
		./railcontrol --logfile=/tmp/railcontrol.log -d
		echo "$!" > /tmp/railcontrol.pid
	}
}

stop() {
	killall railcontrol
}
